# etc/crontabs/root
1 4 * * * /root/CreateOpenAndCloseDoorCronJobs
15 05 * * * /root/AutoDoorControl open
15 21 * * * /root/AutoDoorControl close
